# linkquality_stats.service
[Unit]
Description=WiFi Link Quality Stats Daemon
After=network.target

[Service]
Type=forking
PIDFile=/var/tmp/linkquality_stats.pid
# libmemfnswrap.so hooks allocators via dlsym(RTLD_NEXT,"memalign").
# LD_PRELOADing it places it at link-map position #2 so RTLD_NEXT covers libc.
Environment="LD_PRELOAD=/usr/lib/libmemfnswrap.so.0"
ExecStart=/usr/bin/linkquality_stats
Restart=on-failure

[Install]
WantedBy=multi-user.target
